# srv/salt/node/node-web-skel.conf
[program:node-web-skel]
directory=/opt/src
command=/bin/node app.js
environment=PORT="3000"
numprocs=1
autostart=true
autorestart=true

